# prod.crontab
# Needs to use /bin/bash to make it load our .bash_profile, which sets up the
# $PATH required for uv

*/3 * * * * (cd /home/vitor/Projetos/rpi-media && /bin/bash -lc "make ensure_video_is_playing")
0 12,18 * * * (cd /home/vitor/Projetos/rpi-media && /bin/bash -lc "make play_musica")
0 10,16 * * * (cd /home/vitor/Projetos/rpi-media && /bin/bash -lc "make play_sessao_da_tarde")
0 15 * * * (cd /home/vitor/Projetos/rpi-media && /bin/bash -lc "make play_viagens_brasil")
