# prod.crontab
# Needs to use /bin/bash to make it load our .bash_profile, which sets up the
# $PATH required for uv

*/3 6-19 * * * (cd /home/vitor/Projetos/rpi-media && /bin/bash -lc "make ensure_video_is_playing") 2>&1 | systemd-cat -t rpimedia-cron
*/3 2-7 * * * (cd /home/vitor/Projetos/rpi-media && /bin/bash -lc "make mute_before_dawn") 2>&1 | systemd-cat -t rpimedia-cron
0 11,18 * * * (cd /home/vitor/Projetos/rpi-media && /bin/bash -lc "make play_musica") 2>&1 | systemd-cat -t rpimedia-cron
0 10,16 * * * (cd /home/vitor/Projetos/rpi-media && /bin/bash -lc "make play_sessao_da_tarde") 2>&1 | systemd-cat -t rpimedia-cron
0 15 * * * (cd /home/vitor/Projetos/rpi-media && /bin/bash -lc "make play_viagens_brasil") 2>&1 | systemd-cat -t rpimedia-cron
